# pkgsrc/etc/profile.d/netbsd-pkgsrc.sh
# location of the NetBSD binaries
export PATH=$PATH:/usr/pkg/sbin:/usr/pkg/bin
export MANPATH=$MANPATH:/usr/pkg/man
