# etc/scripts/rc.local
#!/bin/sh

list=$(ls /etc/scripts/startup/*)
for f in $list; do
	if [ -x $f ]; then
		echo "rc.local: $f"
		$f
	fi
done
echo "rc.local: done"
