# etc/cron.daily/fstrim-system
#!/bin/sh
TAG=`basename $0`
#fstrim -v / | logger -t $TAG
#fstrim -v /mnt/home | logger -t $TAG
