# etc/cron.daily/fstrim-system
#!/bin/sh
TAG=`basename $0`
fstrim -av / | logger -t $TAG
