# etc/cron.daily/motd.sh
#!/bin/sh
fortune diogenis >> /etc/motd

