# etc/cron.weekly/clean-tmp
#!/bin/sh
TAG=`basename $0`
find /tmp -type f -ctime +7 -delete
logger -t $TAG "cleanup /tmp"
