# etc/cron.weekly/fstrim-system
#!/bin/sh
TAG=`basename $0`
fstrim -a | logger -t $TAG
