# netbsd-pkgsrc/etc/profile.d/netbsd-pkgsrc.csh
# location of the NetBSD binaries
setenv PKGSRC_PREFIX	/usr/pkgsrc
setenv PATH				${PATH}:/usr/pkg/sbin:/usr/pkg/bin
setenv MANPATH			${MANPATH}:/usr/pkg/man
