# pulseaudio-system/runit/pulseaudio/run
#!/bin/sh
exec pulseaudio --system --realtime=no --high-priority=no --disallow-exit
