# dist/systemd/pulseaudio-system.service
#
#	Use PulseAudio as system-wide daemon
#

[Unit]
Description=PulseAudio Daemon

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
PrivateTmp=true
#ExecStart=/usr/bin/pulseaudio --system --realtime --disallow-exit --no-cpu-limit 
ExecStart=/usr/bin/pulseaudio --system --realtime=no --disallow-exit --enable-memfd=yes --high-priority=no
