# dist/systemd/rc-resume.service
#
#	executes /etc/rc.resume after sleep
#
# systemctl enable rc-resume
# systemctl daemon-reload

[Unit]
Description=Executes /etc/rc.resume after sleep
After=suspend.target

[Service]
Type=oneshot
ExecStart=/etc/rc.resume

[Install]
WantedBy=suspend.target

