# devuan/etc/rc.local
#!/bin/sh -e

# set power-manager to magnetic disks
hdparm -S 36 /dev/sdb
hdparm -S 36 /dev/sdd
