# devuan/etc/rc.local
#!/bin/sh -e

#
ethname=eno1
ip link set $ethname down
ip link set dev $ethname name eth0
ip addr change 192.168.1.100 dev eth0
ip link set eth0 up

# set power-manager to magnetic disks
hdparm -S 36 /dev/sdb
hdparm -S 36 /dev/sdd
