# misc/netbsd-pkgsrc/etc/profile.d/netbsd-pkgsrc.sh
# location of the NetBSD binaries
export PKGSRC_PREFIX=/usr/pkgsrc
export PATH=$PATH:/usr/pkg/sbin:/usr/pkg/bin
export MANPATH=$MANPATH:/usr/pkg/man
